# files/etc/rc.local
#!/bin/sh

echo "net.core.default_qdisc=fq" >> /etc/sysctl.conf
echo "net.ipv4.tcp_congestion_control=bbr" >> /etc/sysctl.conf
echo "vm.swappiness=10" >> /etc/sysctl.conf

exit 0
